# gaia-bigdata/scripts/sv/lweconnectors/run
#!/bin/sh

JVM=`which java`
SVWAIT="25"

exec 2>&1
exec chpst -u hadoop $($JVM -jar /data/lucid/lucidworks/lwe/app/bin/StartUtil.jar -jvm_cmd $JVM -only connectors ${1+"$@"} --print start)






